# home/user/.config/autostart/Finish1.desktop
[Desktop Entry]
Exec=/usr/bin/bash /home/user/Desktop/Finish1.sh
Name=Finish1
Type=Application
Version=1.0
