# home/user/.config/autostart/lxqt-qlipper-autostart.desktop
[Desktop Entry]
Exec=qlipper
Hidden=true
Name=Qlipper
OnlyShowIn=LXQt;
TryExec=qlipper
Type=Application
X-LXQt-Need-Tray=true
