# home/user/.config/autostart/lxqt-xscreensaver-autostart.desktop
[Desktop Entry]
Exec=xscreensaver -no-splash
Hidden=true
Name=XScreenSaver
OnlyShowIn=LXQt;
TryExec=xscreensaver
Type=Application
Version=1.0
